# debian/prerm
#!/bin/bash
set -e

# Stop and disable the service before removal
if [ "$1" = "remove" ] || [ "$1" = "upgrade" ]; then
    systemctl stop star2gps.service || true
    systemctl disable star2gps.service || true
fi

exit 0

